# scripts/cicerone.service
[Unit]
Description=Cicerone Go Autonomous Agent
Documentation=https://github.com/crab-meat-repos/cicerone-goclaw
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=cicerone
Group=cicerone
WorkingDirectory=/opt/cicerone
ExecStart=/opt/cicerone/cicerone server --host 127.0.0.1 --port 18789
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=OLLAMA_HOST=http://127.0.0.1:11434

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/cicerone
PrivateTmp=yes

[Install]
WantedBy=multi-user.target